# deploy/filamind-flow.service
# FilaMind Flow backend — systemd unit (EXAMPLE).
# The recommended installer (scripts/install.sh) generates this file with YOUR real
# user and install path. If installing manually instead: replace every `pi` /
# `/home/pi` below with your user/home (e.g. biqu on a BTT CB1, mks, or any other
# host user), copy to /etc/systemd/system/filamind-flow.service, then:
#   sudo systemctl daemon-reload && sudo systemctl enable --now filamind-flow
[Unit]
Description=FilaMind Flow backend (Klipper / Moonraker panel)
After=network-online.target moonraker.service
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/filamind-flow/backend
EnvironmentFile=-/home/pi/filamind-flow/backend/.env
ExecStart=/home/pi/filamind-flow/backend/.venv/bin/python -m app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
